# rclone/crontab
51 15 * * * /backup.sh >> /share/logs/rclone 2>&1
