# rclone/crontab
40 16 * * * /backup.sh >> /share/logs/rclone 2>&1
