# rclone/crontab
49 15 * * * /backup.sh >> /share/logs/rclone 2>&1
# End this file with an empty new line
Tue Nov 19 03:10:11 PM UTC 2024